# launch/deploy/systemd/oracle-node@.service
# =============================================================================
# launch/deploy/systemd/oracle-node@.service
#
# A templated systemd unit, one instance per cluster node. Usage:
#
#   sudo systemctl enable --now oracle-node@0.service
#   sudo systemctl enable --now oracle-node@1.service
#   ...
#
# The instance ID (%i) is the node's cluster index (0..4 for a 5-node
# cluster). Each instance reads its own env file at
# /etc/phylanx/oracle-node-%i.env.
#
# Hardening:
#   * Runs as a dedicated `phylanx` user.
#   * NoNewPrivileges, ProtectSystem, ProtectHome — process can't escalate.
#   * MemoryMax + CPUQuota bound resource use.
#   * Restart=on-failure with bounded retries — does NOT loop on
#     ProductionRefused (exit code 2) so a misconfigured mainnet start
#     fails LOUDLY rather than thrashing.
# =============================================================================

[Unit]
Description=Phylanx Oracle Node (cluster instance %i)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=phylanx
Group=phylanx
WorkingDirectory=/opt/phylanx

# Per-instance env file — pinned cluster id, port, peer list, network.
EnvironmentFile=/etc/phylanx/oracle-node-%i.env

# The entrypoint refuses to start against mainnet without explicit
# PHYLANX_MAINNET_OK=1 in the env file.
ExecStart=/opt/phylanx/venv/bin/python -m oracle.cluster.run_cluster_node \
    --node-id ${PHYLANX_NODE_ID} \
    --port    ${PHYLANX_PORT} \
    --peers   ${PHYLANX_PEERS}

# Restart policy: bounded retries for transient failures, but NOT for the
# guard-refusal exit code (2). systemd's `RestartPreventExitStatus=`
# blocks restart on those exits.
Restart=on-failure
RestartSec=10s
StartLimitIntervalSec=300s
StartLimitBurst=5
RestartPreventExitStatus=2

# Process hardening.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true

# ── VULN-25 supply-chain hardening ──────────────────────────────────────────
# Read-only root + tightly-scoped writable paths. The venv at
# /opt/phylanx is read-only — a compromised dependency cannot
# silently rewrite its own bytes on disk to survive a restart.
# Writable state is confined to /var/lib/phylanx (db cache, baseline
# snapshots) and /var/log/phylanx (rotated journals).
ReadOnlyPaths=/opt/phylanx
ReadWritePaths=/var/lib/phylanx /var/log/phylanx
ProtectKernelLogs=true
ProtectHostname=true
ProtectClock=true
ProtectProc=invisible
ProcSubset=pid
RestrictRealtime=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources @mount @debug @reboot
SystemCallArchitectures=native
# Drop every capability — the node never needs to listen on a
# privileged port or open a raw socket.
CapabilityBoundingSet=
AmbientCapabilities=
# Network egress: production deploys layer an iptables/nftables
# allowlist in front (see launch/runbooks/supply_chain.md). Systemd
# itself cannot enforce a DNS-based allowlist — IPAddressAllow takes
# IP/CIDR only, set in the env file per-cluster.

# Bound resource use — keeps a runaway node from starving the host.
MemoryMax=2G
CPUQuota=200%
TasksMax=512

# Log to journald with structured rate-limiting.
StandardOutput=journal
StandardError=journal
SyslogIdentifier=phylanx-oracle-%i

[Install]
WantedBy=multi-user.target
